# etf-advisor/deploy/admin-server.service
[Unit]
Description=ETF Advisor Admin API Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/etf-advisor
ExecStart=/opt/etf-advisor/venv/bin/python /opt/etf-advisor/scripts/admin_server.py
Restart=on-failure
RestartSec=5
StandardOutput=append:/opt/etf-advisor/data/admin.log
StandardError=append:/opt/etf-advisor/data/admin.log

[Install]
WantedBy=multi-user.target
